# infra/xcresult-processor-image/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Boots the Tuist server release in xcresult-processor mode. Sourced env
  comes from /etc/tuist.env, populated by /opt/tuist/inject-env.sh at
  first boot from the tart-kubelet env mount
  (/Volumes/My Shared Files/env/tuist.env).

  We exec a wrapper shell rather than calling /opt/tuist/release/bin/tuist
  directly so we can source the env file before launching. launchd's
  EnvironmentVariables key can't reference a file at load time — but it
  does take literal values, which we use for SECRETS_DIRECTORY.

  SECRETS_DIRECTORY points at the release's bundled priv/secrets/.
  `mix release` packages the app's priv/ at lib/<app>-<vsn>/priv/, so a
  relative `priv/secrets/<env>.yml.enc` read against the release cwd
  (/opt/tuist/release) would `enoent`. The Linux Dockerfile sidesteps
  this by copying the encrypted YAMLs to a top-level /app/priv/secrets/
  and pointing SECRETS_DIRECTORY there; the Tart image keeps the files
  where mix put them and just sets the absolute path. Without this the
  BEAM crashloops on boot inside Tuist.Environment.decrypt_secrets/0.

  KeepAlive=true so the BEAM gets restarted on unexpected exit. Liveness
  surfaces through the launchd job staying up + Oban heartbeats into
  oban_peers — there's no HTTP probe (the release runs with TUIST_WEB=0).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.tuist.xcresult-processor</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SECRETS_DIRECTORY</key>
        <string>/opt/tuist/release/lib/tuist-0.1.0/priv/secrets/</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/opt/tuist/inject-env.sh &amp;&amp; source /etc/tuist.env &amp;&amp; exec /opt/tuist/release/bin/tuist start</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/opt/tuist/release</string>

    <key>StandardOutPath</key>
    <string>/var/log/xcresult-processor/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/xcresult-processor/stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
