# infra/runner-image/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Boots the Tuist runner image. The wrapper:
    1. Runs /opt/tuist/inject-env.sh — materializes /etc/tuist.env
       from the tart-kubelet env mount.
    2. Sources /etc/tuist.env so the dispatch-poll script sees
       the per-Pod TUIST_RUNNER_* vars.
    3. Execs /opt/tuist/dispatch-poll.sh — polls the Tuist server
       for a JIT config and execs ./run.sh on first 200.

  This is a LaunchAgent loaded from
  /Users/runner/Library/LaunchAgents/. The Packer build configures
  runner auto-login + passwordless sudo, so the VM boots into a
  real desktop session for the `runner` user (matching GitHub-
  hosted's runtime account), loginwindow loads this agent, and
  the agent process inherits the user session.

  Why an agent and not a daemon: LaunchDaemons start with a
  near-empty environment (no HOME, no PATH past launchd's
  defaults, locale = C) and the Actions runner forwards its env
  to every step shell. Customer workflows that read $HOME or
  invoke a Homebrew-installed tool break before the first step
  runs. Wrapping the entrypoint in `/bin/zsh -lc` here makes it
  a login shell, which sources the cirruslabs base's ~/.zprofile
  (Homebrew shellenv, rbenv init, LANG=en_US.UTF-8, PATH
  additions for Android SDK + Flutter + openjdk + Node). Step
  shells inherit the same environment an interactive SSH session
  on the same VM would see. Supersedes the per-var workarounds
  layered on top of the prior LaunchDaemon (#10797, #10800).

  KeepAlive=false: the runner image is single-shot. After the
  GitHub Actions job finishes, dispatch-poll.sh halts the VM so
  `tart run` returns and tart-kubelet transitions the Pod to
  Succeeded. Restarting would just re-poll for a JIT config
  that's already been claimed.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.tuist.runner</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <!--
          inject-env.sh runs under sudo. tart-kubelet stages the
          shared env directory at /Volumes/My Shared Files/env/ as
          0700 root, with tuist.env / sa_token inside at 0600 root.
          The agent runs as the `runner` user, so without elevation
          runner can't even stat the directory: inject-env.sh's
          `[ ! -f ENV_SRC ]` test returns true, the script exits 0
          ("no env source"), and dispatch-poll.sh boots with an
          empty /etc/tuist.env (no TUIST_RUNNER_DISPATCH_URL, no
          SA token). Elevate only the injector; dispatch-poll runs
          as runner so the GitHub runner stays non-root.
        -->
        <string>sudo /opt/tuist/inject-env.sh &amp;&amp; source /etc/tuist.env &amp;&amp; exec /opt/tuist/dispatch-poll.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>/Users/runner/actions-runner</string>

    <key>StandardOutPath</key>
    <string>/var/log/tuist-runner/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/tuist-runner/stderr.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
