# infra/tart-cri/platform/dev.tuist.kubelet.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Boots kubelet on the Mac mini. Dependencies:
    - tart-cri must be running (LaunchDaemon ordering can't express
      this directly, but kubelet retries the CRI socket connection
      on its own).
    - /etc/kubernetes/kubelet-config.yaml + /etc/kubernetes/kubelet.conf
      (kubeconfig) must exist — provisioned by the bootstrap script.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.tuist.kubelet</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/kubelet</string>
        <string>--config=/etc/kubernetes/kubelet-config.yaml</string>
        <string>--kubeconfig=/etc/kubernetes/kubelet.conf</string>
        <string>--bootstrap-kubeconfig=/etc/kubernetes/bootstrap-kubelet.conf</string>
        <string>--cert-dir=/var/lib/kubelet/pki</string>
        <string>--hostname-override=$(hostname)</string>
        <string>--node-labels=kubernetes.io/os=darwin,kubernetes.io/arch=arm64,tuist.dev/runtime=tart</string>
        <string>--register-with-taints=tuist.dev/macos=true:NoSchedule</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/kubelet/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/kubelet/stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
